# .config/systemd/user/my-login-commands.service
[Unit]
Description=Siva's custom login commands

[Service]
ExecStart=%h/bin/my-login-commands

[Install]
WantedBy=default.target
